# contrib/upstart/captainhook.conf
description "captainhook webhook server"
console log
start on started
stop on shutdown

exec /opt/captainhook/bin/captainhook -configdir /opt/captainhook -listen-addr 0.0.0.0:8080 -echo
